# datacenter/scripts/weather-subscriber.service
[Unit]
Description=Weather Network MQTT Subscriber
After=network-online.target mariadb.service mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=weather
Group=weather
WorkingDirectory=/opt/weather
ExecStart=/opt/weather/venv/bin/python3 /opt/weather/datacenter/mqtt_subscriber.py \
    --config /etc/weather/config.ini
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
